# etc/systemd/system/bootc-fetch-apply-updates.timer.d/weekly.conf
[Timer]
RandomizedDelaySec=600
#Apply updates on Saturday at midnight
OnCalendar=Sat *-*-* 03:00:00

